# relayservice.conf
description "Service for relay_toggle"
author "your@email-address.com"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /home/pi/relay_toggle.py
